# com.user.speech_transcriber.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.speech_transcriber</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Users/dixon/code/whisper/speech_transcriber/venv/bin/python</string>
        <string>/Users/dixon/code/whisper/speech_transcriber/main.py</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardErrorPath</key>
    <string>/Users/dixon/code/whisper/speech_transcriber/error.log</string>
    <key>StandardOutPath</key>
    <string>/Users/dixon/code/whisper/speech_transcriber/output.log</string>
</dict>
</plist>